# deploy/dotlearn-api.service
[Unit]
Description=DotLearn API (NestJS)
Documentation=https://github.com/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_USER__
WorkingDirectory=__API_RUNTIME__
EnvironmentFile=__ENV_FILE__
ExecStart=__NODE_BIN__ dist/main.js
Restart=on-failure
RestartSec=5
SyslogIdentifier=dotlearn-api

# Sandboxing / hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=false
ReadWritePaths=__DATA_DIR__
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
PrivateDevices=true
ProtectKernelLogs=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RemoveIPC=true
UMask=0077
MemoryMax=512M
TasksMax=256

[Install]
WantedBy=multi-user.target
